# charm/trusty/saikuanalytics/templates/upstart.conf
description "saiku"
author "Tom Barber <Tom.Barber@bugg-pc.lan>"

start on runlevel [2345]
stop on runlevel [016]

respawn

console log
script
    echo Fake service; sleeping for an hour...
    sleep 360
end script
